# packages/sing-box-config/sing_box_config-0.4.0-py3-none-any.whl/roles/sing_box_config/templates/etc/systemd/system/sing-box-config-updater.service.j2
[Unit]
Description=sing-box-config updater

[Service]
Type=oneshot
User={{ proxy_user.name }}
Group={{ proxy_group.name }}
WorkingDirectory={{ sing_box_state_dir }}
ExecStart={{ sing_box_venv_dir }}/bin/sing-box-config --base {{ sing_box_config_dir }}/base.json --subscriptions {{ sing_box_config_dir }}/subscriptions.json --output {{ sing_box_etc_dir }}/config.json
